# etc/systemd/system/unit-status-telegram@.service
[Unit]
Description=Unit Status Telegram Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/unit-status-telegram %I